# src/provision/files/vinz.conf
[program:vinz]
command=/usr/bin/python app.py
directory=/vagrant/src/
user=nobody
autostart=true
autorestart=true
