# provision/files/vinz_supervisor.conf
[program:vinz]
command=/usr/bin/python main.py
directory=/vagrant/app/
user=nobody
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/vinz-stdout.log
